# scripts/launchd/com.codex.swing.kr-daily-auto-scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.codex.swing.kr-daily-auto-scan</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/dongdong/Projects/codex_swing/swing-main/scripts/run_kr_daily_auto_scans.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>20</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/dongdong/Projects/codex_swing/swing-main/runtime_state/discord_jobs/kr_daily_auto_scan.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/dongdong/Projects/codex_swing/swing-main/runtime_state/discord_jobs/kr_daily_auto_scan.err.log</string>
</dict>
</plist>
