# scripts/launchd/com.codex.swing.kr-premarket-theme-prior.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.codex.swing.kr-premarket-theme-prior</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/dongdong/Projects/codex_swing/swing-main/scripts/run_kr_daily_auto_scans.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>AG_KR_DAILY_PHASE</key>
    <string>premarket</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>20</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/dongdong/Projects/codex_swing/swing-main/runtime_state/discord_jobs/kr_premarket_theme_prior.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/dongdong/Projects/codex_swing/swing-main/runtime_state/discord_jobs/kr_premarket_theme_prior.err.log</string>
</dict>
</plist>
